# scripts/launchd/com.user.claude-weekly-lints.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.user.claude-weekly-lints  -- weekly off-hours wiki-lint + skill-lint sweep.
  Req 1 of coo-cockpit-upgrades. This is the repo-tracked source of truth; the live
  copy is installed to ~/Library/LaunchAgents/ and loaded via launchctl.
  Reinstall:  cp this file ~/Library/LaunchAgents/ && launchctl load -w ~/Library/LaunchAgents/com.user.claude-weekly-lints.plist
  Fire now:   launchctl kickstart gui/$(id -u)/com.user.claude-weekly-lints
  Monitor:    tail -f ~/.claude/state/lint-cron.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.claude-weekly-lints</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/you/.claude/scripts/run-weekly-lints.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>21</integer>
    <key>Minute</key><integer>7</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/you/.claude/state/lint-cron.launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/you/.claude/state/lint-cron.launchd.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
